# src/alise/alise.path
[Path]
PathExists=/dev/hidraw0
Unit=alise.service 
