# src/alise/alise.service.in
[Unit]
Description=alise

[Service]
Type=simple
ExecStart=@EXE_PATH@

[Install]
WantedBy=multi-user.target
