# src/alise/alise.service.in
[Unit]
Description=alise

[Service]
Type=simple
ExecStart=@EXE_PATH@
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
